# root/etc/init.d/pw2watchdog
#!/bin/sh /etc/rc.common
# PW2WD_VERSION: v0.4.0-dev

START=98
STOP=15

WATCHDOG=/usr/bin/pw2watchdog.sh
SCANNER=/usr/bin/pw2watchdog-scanner.sh
SUBSCRIBE=/usr/bin/pw2watchdog-subscribe.sh

PID_WATCHDOG=/var/run/pw2watchdog/watchdog.pid
PID_SCANNER=/var/run/pw2watchdog/scanner.pid

ENV_FILE=/var/run/pw2watchdog/env.static

# Cleanup stale blackhole nft rule using env-resolved table/chain (Layered Robustness:
# fall back to env.static so this works on any PassWall2 fork / non-standard naming).
# If env.static is absent (first boot, env resolver not yet run) — skip cleanup with
# a logger note instead of guessing.
cleanup_stale_blackhole() {
	local phase="$1"  # "start" | "stop"
	local state_file="/var/run/pw2watchdog/state"
	[ -f "$state_file" ] || return 0

	local handle
	handle="$(grep 'STATIC_BH_HANDLE' "$state_file" | cut -d"'" -f2)"
	[ -n "$handle" ] && [ "$handle" != "unknown" ] || return 0

	if [ ! -f "$ENV_FILE" ]; then
		logger -t pw2watchdog "$phase: env.static missing, skipping nft cleanup (handle=$handle)"
		return 0
	fi

	# shellcheck disable=SC1090
	. "$ENV_FILE"
	local tbl="$PW2_NFTABLE_NAME"
	local chn="$PW2_NFTCHAIN_MANGLE"
	if [ -z "$tbl" ] || [ -z "$chn" ]; then
		logger -t pw2watchdog "$phase: nft table/chain not resolved in env.static, skipping cleanup (handle=$handle)"
		return 0
	fi

	# IMPORTANT: $tbl deliberately unquoted (word-splits to family+name, e.g. 'inet passwall2').
	# This matches the convention used throughout pw2watchdog.sh.
	# shellcheck disable=SC2086
	nft delete rule $tbl "$chn" handle "$handle" 2>/dev/null && \
		logger -t pw2watchdog "$phase: removed stale blackhole nft handle $handle from $tbl/$chn"
}

start() {
	# Remove any stale procd registration to prevent ubus conflicts
	ubus call service delete '{"name":"pw2watchdog"}' 2>/dev/null

	local enabled
	enabled="$(uci -q get pw2watchdog.main.enabled)"
	[ "$enabled" = "1" ] || {
		logger -t pw2watchdog "Service is disabled, not starting."
		return 0
	}

	mkdir -p /var/run/pw2watchdog

	# Clean up stale blackhole rule from previous run (Layered Robustness: use env.static)
	cleanup_stale_blackhole start

	# Subscribe install (idempotent)
	if [ -x "$SUBSCRIBE" ]; then
		"$SUBSCRIBE" install
		local on_boot boot_delay
		on_boot="$(uci -q get pw2watchdog.advanced.sub_update_on_boot)"
		if [ "$on_boot" = "1" ]; then
			boot_delay="$(uci -q get pw2watchdog.advanced.sub_boot_delay)"
			# Minimum delay 120s — ensures scanner completes first cycle before subscribe
			case "$boot_delay" in
				''|*[!0-9]*) boot_delay=120 ;;
				*) [ "$boot_delay" -ge 120 ] || boot_delay=120 ;;
			esac
			logger -t pw2watchdog "sub_update_on_boot=1, subscription update delayed ${boot_delay}s"
			( sleep "$boot_delay"; "$SUBSCRIBE" run ) &
		fi
	fi

	# Start watchdog daemon
	if [ -x "$WATCHDOG" ]; then
		"$WATCHDOG" daemon &
		echo $! > "$PID_WATCHDOG"
		logger -t pw2watchdog "watchdog started (pid $!)"
	fi

	# Start scanner daemon
	if [ -x "$SCANNER" ]; then
		"$SCANNER" daemon &
		echo $! > "$PID_SCANNER"
		logger -t pw2watchdog "scanner started (pid $!)"
	fi
}

stop() {
	logger -t pw2watchdog "Stopping service."

	# Clean up static blackhole nft rule (Layered Robustness: use env.static)
	cleanup_stale_blackhole stop

	# Kill all instances (by name — covers manually started processes too)
	killall pw2watchdog.sh pw2watchdog-scanner.sh 2>/dev/null
	sleep 2
	# Force kill if still running
	killall -9 pw2watchdog.sh pw2watchdog-scanner.sh 2>/dev/null
	# Clean up pid files
	rm -f "$PID_WATCHDOG" "$PID_SCANNER"
}

restart() {
	stop
	sleep 1
	start
}
