# root/etc/init.d/pw2watchdog
#!/bin/sh /etc/rc.common

START=98
STOP=15

WATCHDOG=/usr/bin/pw2watchdog.sh
SCANNER=/usr/bin/pw2watchdog-scanner.sh
SUBSCRIBE=/usr/bin/pw2watchdog-subscribe.sh

PID_WATCHDOG=/var/run/pw2watchdog/watchdog.pid
PID_SCANNER=/var/run/pw2watchdog/scanner.pid

start() {
	# Remove any stale procd registration to prevent ubus conflicts
	ubus call service delete '{"name":"pw2watchdog"}' 2>/dev/null

	local enabled
	enabled="$(uci -q get pw2watchdog.main.enabled)"
	[ "$enabled" = "1" ] || {
		logger -t pw2watchdog "Service is disabled, not starting."
		return 0
	}

	mkdir -p /var/run/pw2watchdog

	# Clean up stale blackhole rule from previous run
	local state_file="/var/run/pw2watchdog/state"
	if [ -f "$state_file" ]; then
		local handle
		handle="$(grep 'STATIC_BH_HANDLE' "$state_file" | cut -d"'" -f2)"
		if [ -n "$handle" ] && [ "$handle" != "unknown" ] && [ "$handle" != "" ]; then
			nft delete rule inet passwall2 PSW2_MANGLE handle "$handle" 2>/dev/null && \
				logger -t pw2watchdog "start: removed stale blackhole nft handle $handle"
		fi
	fi

	# Subscribe install (idempotent)
	if [ -x "$SUBSCRIBE" ]; then
		"$SUBSCRIBE" install
		local on_boot boot_delay
		on_boot="$(uci -q get pw2watchdog.advanced.sub_update_on_boot)"
		if [ "$on_boot" = "1" ]; then
			boot_delay="$(uci -q get pw2watchdog.advanced.sub_boot_delay)"
			# Minimum delay 120s — ensures scanner completes first cycle before subscribe
			[ -n "$boot_delay" ] && [ "$boot_delay" -ge 120 ] 2>/dev/null || boot_delay=120
			logger -t pw2watchdog "sub_update_on_boot=1, subscription update delayed ${boot_delay}s"
			( sleep "$boot_delay"; "$SUBSCRIBE" run ) &
		fi
	fi

	# Start watchdog daemon
	if [ -x "$WATCHDOG" ]; then
		"$WATCHDOG" daemon &
		echo $! > "$PID_WATCHDOG"
		logger -t pw2watchdog "watchdog started (pid $!)"
	fi

	# Start scanner daemon
	if [ -x "$SCANNER" ]; then
		"$SCANNER" daemon &
		echo $! > "$PID_SCANNER"
		logger -t pw2watchdog "scanner started (pid $!)"
	fi
}

stop() {
	logger -t pw2watchdog "Stopping service."

	# Clean up static blackhole nft rule
	local state_file="/var/run/pw2watchdog/state"
	if [ -f "$state_file" ]; then
		local handle
		handle="$(grep 'STATIC_BH_HANDLE' "$state_file" | cut -d"'" -f2)"
		if [ -n "$handle" ] && [ "$handle" != "unknown" ] && [ "$handle" != "" ]; then
			nft delete rule inet passwall2 PSW2_MANGLE handle "$handle" 2>/dev/null && \
				logger -t pw2watchdog "stop: removed blackhole nft handle $handle"
		fi
	fi

	# Kill all instances (by name — covers manually started processes too)
	killall pw2watchdog.sh pw2watchdog-scanner.sh 2>/dev/null
	sleep 2
	# Force kill if still running
	killall -9 pw2watchdog.sh pw2watchdog-scanner.sh 2>/dev/null
	# Clean up pid files
	rm -f "$PID_WATCHDOG" "$PID_SCANNER"
}

restart() {
	stop
	sleep 1
	start
}
